# launchd/com.nischal.claudebot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nischal.claudebot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/nischal/bin/claude-bot-start.sh</string>
    </array>

    <!-- Environment — launchd does NOT inherit your shell env -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/nischal/.local/bin:/Users/nischal/.npm-global/bin:/Users/nischal/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/nischal</string>
        <!-- ANTHROPIC_API_KEY is loaded from ~/.claude/.env (chmod 600) by the wrapper script -->
    </dict>

    <!-- Start on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart on exit -->
    <key>KeepAlive</key>
    <true/>

    <!-- Wait 10s between restarts to prevent crash loops -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Logs -->
    <key>StandardOutPath</key>
    <string>/Users/nischal/.claude/logs/claudebot-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/nischal/.claude/logs/claudebot-launchd.log</string>
    <key>SizeLimit</key>
    <integer>52428800</integer>
</dict>
</plist>
