# launchd/com.nischal.claudebot-companion.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nischal.claudebot-companion</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/nischal/.bun/bin/bun</string>
    <string>run</string>
    <string>/Users/nischal/projects/claude-telegram-bot/companion/src/index.ts</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/nischal</string>
    <key>PATH</key>
    <string>/Users/nischal/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/nischal/projects/claude-telegram-bot/companion</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/nischal/.claude/logs/companion.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/nischal/.claude/logs/companion.log</string>
</dict>
</plist>
